# src/target-rootfs/etc/systemd/system/blade-firstboot-gpu.service
[Unit]
Description=Validate hybrid graphics before enabling GNOME
Wants=switcheroo-control.service
After=systemd-modules-load.service switcheroo-control.service
ConditionPathExists=!/var/lib/blade-installer/gpu-gate-passed
ConditionPathExists=!/var/lib/blade-installer/gpu-gate-failed

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/blade-firstboot-gpu
StandardOutput=journal+console
StandardError=journal+console

[Install]
WantedBy=multi-user.target
